# etc/crontabs/root
* * * * * cd /root/mqtt_client/ ; ./mqtt_send.lua "sonny_lua_time" "time" $(date +%s) 1> /dev/null
* * * * * cd /root/mqtt_client/ ; ./mqtt_send.lua "sonny_lua_time" "time/local" $(./mqtt_time.sh) 1> /dev/null
