# etc/crontabs/root
* * * * * cd /root/mqtt_client/ ; ./mqtt_send.lua "sonny_lua_time" "time" $(date +%s) 2 0> /dev/null